# systemd/openli-mediator.service
[Unit]
Description=OpenLI mediator daemon
Documentation=http://github.com/wanduow/openli/wiki

[Service]
Type=simple
ExecStart=/usr/bin/openlimediator -c /etc/openli/mediator-config.yaml
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
